# recipes-tools/probe-rs/files/probe-rs-dap-server.service
[Unit]
Description=probe-rs cmd server 
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/probe-rs dap-server --ip 0.0.0.0 --port 50000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target